# ilovemoney/conf-templates/supervisord.conf.j2
[program:ilovemoney]
command={{ bin_path }}/gunicorn -c /etc/ilovemoney/gunicorn.conf.py ilovemoney.wsgi:application
user=ilovemoney
autostart=true
autorestart=true
redirect_stderr=true
